# packaging/systemd/kbd-auto-layout.service
[Unit]
Description=Automatic keyboard layout switching based on connected keyboards
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kbd-auto-layoutd
Restart=always
RestartSec=2
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority

[Install]
WantedBy=default.target
